# deploy/systemd/whs-nhs-hours.service
[Unit]
Description=WHS NHS Hours (Next.js)
After=network.target
# If startup (build) fails 3x within 10 min, stop retrying instead of rebuild-looping.
StartLimitIntervalSec=600
StartLimitBurst=3

[Service]
Type=simple
User=musicserver
Group=musicserver
WorkingDirectory=/home/musicserver/whs-nhs-hours
EnvironmentFile=/home/musicserver/whs-nhs-hours/.env
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
# If Node is installed via nvm, include its bin path.
Environment=PATH=/home/musicserver/.nvm/versions/node/v24.17.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin

# Deploy flow: `git pull` (then `npm ci` if deps changed) and restart the service.
# prestart.sh applies migrations and rebuilds ONLY when HEAD changed since the last
# build (otherwise it reuses the old .next, so an unchanged restart is fast).
ExecStartPre=/usr/bin/env bash deploy/prestart.sh
ExecStart=/usr/bin/env npm run start
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
# The build runs during startup; give it room on a small box.
TimeoutStartSec=300

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
# Node/V8 JIT needs RW->RX memory permission changes.
MemoryDenyWriteExecute=false
SystemCallArchitectures=native
# Next.js startup queries local interfaces (uv_interface_addresses), which needs AF_NETLINK.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
# The ops console edits files and runs git (fetch/pull), so the whole project dir
# must be writable. The rest of the system stays read-only (ProtectSystem=strict)
# and privilege escalation stays blocked (NoNewPrivileges) — this only widens writes
# to the app's own directory, nothing else.
ReadWritePaths=/home/musicserver/whs-nhs-hours

[Install]
WantedBy=multi-user.target
